# certbot/crontab
# Single container cron example
* * * * * cert-renew-script.sh  >> /etc/letsencrypt/foo 2>&1
# -> Every 12 hours