# certbot/crontab
# Single container cron example
* * * * * ./cert-renew-script.sh  >> /root/logs/certbot_cron 2>&1
# -> Every 12 hours