# certbot/crontab
# Single container cron example
0 */12 * * * root bash cert-renew-script.sh
# -> Every 12 hours