# certbot/crontab
# Single container cron example
* * * * * root bash cert-renew-script.sh
# -> Every 12 hours